# modules/supervisord/files/supervisord.d/keepalived.ini
[program:keepalived]
command=/usr/sbin/keepalived -D -n
stopsignal=TERM
priority=30
autostart=false
