# modules/supervisord/files/supervisord.d/redis.ini
[program:redis6379]
command=/usr/sbin/redis-server /etc/redis/redis-6379.conf
stopsignal=TERM
priority=50
autorestart=true
startsecs=10
user=redis
stdout_logfile=/var/log/redis/redis-6379.log
autostart=false

[program:redis6380]
command=/usr/sbin/redis-server /etc/redis/redis-6380.conf
stopsignal=TERM
priority=50
autorestart=true
startsecs=10
user=redis
stdout_logfile=/var/log/redis/redis-6380.log
autostart=false

[program:redis6381]
command=/usr/sbin/redis-server /etc/redis/redis-6381.conf
stopsignal=TERM
priority=50
autorestart=true
startsecs=10
user=redis
stdout_logfile=/var/log/redis/redis-6381.log
autostart=false

[program:redis6383]
command=/usr/sbin/redis-server /etc/redis/redis-6383.conf
stopsignal=TERM
priority=50
autorestart=true
startsecs=10
user=redis
stdout_logfile=/var/log/redis/redis-6383.log
autostart=false

[program:redis6385]
command=/usr/sbin/redis-server /etc/redis/redis-6385.conf
stopsignal=TERM
priority=50
autorestart=true
startsecs=10
user=redis
stdout_logfile=/var/log/redis/redis-6385.log
autostart=false

[program:redis6386]
command=/usr/sbin/redis-server /etc/redis/redis-6386.conf
stopsignal=TERM
priority=50
autorestart=true
startsecs=10
user=redis
stdout_logfile=/var/log/redis/redis-6386.log
autostart=false

[program:redis6387]
command=/usr/sbin/redis-server /etc/redis/redis-6387.conf
stopsignal=TERM
priority=50
autorestart=true
startsecs=10
user=redis
stdout_logfile=/var/log/redis/redis-6387.log
autostart=false

[program:redis6388]
command=/usr/sbin/redis-server /etc/redis/redis-6388.conf
stopsignal=TERM
priority=50
autorestart=true
startsecs=10
user=redis
stdout_logfile=/var/log/redis/redis-6388.log
autostart=false
